# catch-a-tweet.service
[Unit]
Description=catch-a-tweet — X account watcher and Telegram alerter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/libs/catch-a-tweet
Environment=PYTHONPATH=/home/azureuser/libs/catch-a-tweet/src
EnvironmentFile=/home/azureuser/libs/catch-a-tweet/.env
ExecStart=/home/azureuser/libs/catch-a-tweet/.venv/bin/python src/runner.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
